# clean-slack.cron
*/15 * * * * cd /opt/clean-slack/ && python /opt/clean-slack/clean-slack.py
