# clean-slack.cron
*/20 * * * * cd /opt/clean-slack/ && python /opt/clean-slack/clean-slack.py
*/18 * * * * cd /opt/clean-slack/ && python /opt/clean-slack/clean-slack-channel.py
