# clean-slack.cron
*/10 * * * * cd /opt/clean-slack/ && python /opt/clean-slack/clean-slack.py
*/6 * * * * cd /opt/clean-slack/ && python /opt/clean-slack/clean-slack-channel.py
