# clean-slack.cron
0 * * * * /opt/clean-slack/clean-slack.py
