# clean-slack.cron
*/2 * * * * cd /opt/clean-slack/ && python /opt/clean-slack/clean-slack.py
*/3 * * * * cd /opt/clean-slack/ && python /opt/clean-slack/clean-slack-channel.py
